# src/main/resources/spark-submit-supervisor.conf
[program:spark-submit]
command=sh /opt/apache-spark/sbin/spark-submit-app.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/stdout.log
autostart=true
autorestart=false